# docker/trello2mail.cron
# Private tasks
30 10 * * * TRELLO_URL="" /app/trello2mail

# Project 1
30 10 * * * TRELLO_URL="" /app/trello2mail

# Project 2
30 10 * * * TRELLO_URL="" /app/trello2mail
