# docker/trello2mail.cron
# Private tasks
0 10 * * * TRELLO_URL="" /app/trello2mail

# Project 1
0 10 * * * TRELLO_URL="" /app/trello2mail

# Project 2
0 10 * * * TRELLO_URL="" /app/trello2mail
